# coolify-mcp.service
[Unit]
Description=Coolify MCP Server
Documentation=https://github.com/wrediam/coolify-mcp-server
After=network.target

[Service]
Type=simple
User=avi
Group=avi
WorkingDirectory=/home/avi/projects/coolify/coolify-mcp
Environment="NODE_ENV=production"
Environment="COOLIFY_BASE_URL=https://coolify.theprofitplatform.com.au"
Environment="COOLIFY_TOKEN=***REMOVED***"

# Run the built MCP server
ExecStart=/usr/bin/node /home/avi/projects/coolify/coolify-mcp/build/index.js

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=coolify-mcp

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/avi/projects/coolify/coolify-mcp

[Install]
WantedBy=multi-user.target
